# deploy/cestacontrol.service
[Unit]
Description=CestaControl FastAPI
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/cestacontrol
EnvironmentFile=/var/www/cestacontrol/.env
ExecStart=/var/www/cestacontrol/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
